# files/supervisor/conf.d/crond.conf
[program:crond]
command=/usr/sbin/cron -f
process_name=%(program_name)s
autostart=false
autorestart=true
user=root
directory=/
stdout_logfile=/var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s/stderr.log
