# files/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
process_name=%(program_name)s
autostart=true
autorestart=true
user=root
directory=/
stdout_logfile=/var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s/stderr.log